# puppet/precip/templates/restart_services_once_mounted.conf.erb
# restart services once mounted
start on vagrant-mounted

script
  sudo service apache2 restart
  sudo service mysql restart
end script